# systemd/raspilapse-daily-video.timer
[Unit]
Description=Daily Timelapse Video Generation Timer
# Note: Removed "Requires=" - timers trigger their service automatically by name
# Having Requires= caused the service to start immediately when the timer started on boot

[Timer]
# Run daily at 05:00
OnCalendar=*-*-* 05:00:00

# Disabled to prevent duplicate runs after reboot
# (If Pi reboots after 05:00, Persistent=true would trigger another run)
Persistent=false

# Randomize the start time by up to 5 minutes to avoid load spikes
RandomizedDelaySec=300

[Install]
WantedBy=timers.target